# docker/s6-rc.d/wg-obfuscator/finish
#!/command/with-contenv sh
pkill -f '/usr/local/bin/wg-obfuscator' 2>/dev/null
exit 0
