# docker/s6-rc.d/wg-obfuscator/run
#!/command/with-contenv sh

while [ ! -s /run/wg-obfuscator.conf ]; do
  sleep 1
done

exec setsid /usr/local/bin/wg-obfuscator --config /run/wg-obfuscator.conf
