# docker/s6-rc.d/node/run
#!/command/with-contenv sh
cd /app
if [ "${TLS_TERMINATION:-internal}" != "external" ] && [ -f /app/certs/active/fullchain.pem ] && [ -f /app/certs/active/privkey.pem ]; then
  export NITRO_SSL_CERT=$(cat /app/certs/active/fullchain.pem)
  export NITRO_SSL_KEY=$(cat /app/certs/active/privkey.pem)
fi
exec node server/index.mjs
